# aroio/overlay/usr/lib/systemd/system/bluetooth-dbbackup.service
[Unit]
Description=Bluetooth database backup
After=sys-subsystem-bluetooth-devices-hci0.device
ConditionPathIsDirectory=/sys/class/bluetooth/hci0
Before=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bluetooth-dbbackup

[Install]
WantedBy=multi-user.target
