# aroio/overlay/lib/systemd/system/lircd.service
[Unit]
Description=Flexible IR remote input/output application support
After=network.target

[Service]
Type=forking
PIDFile=/var/run/lirc/lircd.pid
ExecStart=/usr/sbin/lircd

[Install]
WantedBy=multi-user.target
