# aroio/overlay/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/bin/infoscreen
StandardInput=tty
StandardOutput=tty
